# ffmpeg-bot.conf
description "Mediacore ffmpeg-bot"
author "Valerio Minetti"

start on startup
stop on shutdown
respawn

exec sudo -u www-data php -f /opt/ffmpeg-bot.php
